# .config/systemd/user/hibernate.service
[Unit]
Description=Hibernate system when battery drops below 5%

[Service]
Type=oneshot
ExecStart=/home/sarius/.hibernateOnLowBat.sh
